# docs/source/setup/supervisor-uwsgi.conf
[program:uwsgi]
command = /path/to/uwsgi --emperor /etc/uwsgi/apps
stdout_logfile = syslog
stderr_logfile = syslog
startsecs = 5
stopsignal = QUIT
stopasgroup = true
killasgroup = true
